# config_files/auto_init
#!/sbin/openrc-run
command=/sbin/auto_init.sh